# systemd/tell-agent.service
[Unit]
Description=Tell Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tell-agent --config /etc/tell/agent.toml
Restart=always
RestartSec=5
User=tell-agent
Group=tell-agent
# Allow reading /proc, /sys, log files
ReadOnlyPaths=/proc /sys /var/log
StateDirectory=tell-agent

[Install]
WantedBy=multi-user.target
